# bikeshed.conf
[program:bikeshed]
priority=200
directory=/BikeShed
command=rails s
autorestart=false
